# vendor-server/systemd/erp-data.service.d/10-restart-vendor-api.conf
# V-05: khi sshfs remount (reboot / rớt mạng dài / crash+restart), container
# sc-vendor-api giữ SLAVE của mount FUSE CŨ (prop=rprivate, bind vào SUBDIR của
# mountpoint) → thấy thư mục RỖNG im lặng dù host thấy đủ file (đã chứng minh T3).
# rslave KHÔNG fix được (remount là event trên parent mount, không propagate tới
# slave của mount cũ). Fix deterministic: sau mỗi lần (re)mount sshfs xong → chờ
# mountpoint sẵn sàng rồi restart container để nó bind lại FUSE mới.
# Dấu '-' để lúc boot (docker chưa lên) ExecStartPost fail không làm hỏng unit.
[Service]
ExecStartPost=-/bin/sh -c 'for i in $(seq 1 30); do mountpoint -q /mnt/erp-data && break; sleep 1; done; docker restart sc-vendor-api 2>/dev/null || true'
